# crontab
0 20 * * 1 root /usr/local/bin/node script.js >> /var/log/cron.log 2>&1